# xinetd.d/irctelegram
service irctelegram
{
    disable         = no
    user            = irctelegram
    wait            = no
    server          = /etc/irctelegram/irctelegram.py
    socket_type     = stream
    protocol        = tcp
    port            = 28281
    type            = UNLISTED
    log_on_success  = HOST PID USERID DURATION EXIT
}
